# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run the ESPHome log viewer service
# ==============================================================================

# Declare variables
declare device_ip
declare password
declare encryption_key
declare port

# Get configuration values
device_ip=$(bashio::config 'device_ip')
password=$(bashio::config 'password')
encryption_key=$(bashio::config 'encryption_key')
port=$(bashio::config 'port')

# Check if required fields are provided
if [ -z "$device_ip" ] || [ -z "$port" ]; then
  bashio::log.error "Device IP and port are required!"
  exit 1
fi

bashio::log.info "Connecting to ESPHome device at ${device_ip}:${port}..."

# Export variables to environment for Python
export DEVICE_IP="${device_ip}"
export PASSWORD="${password}"
export ENCRYPTION_KEY="${encryption_key}"
export PORT="${port}"

# Execute the Python script with environment variables
exec python3 /tmp/esphome_log.py